# etc/service/apache/run
#!/bin/bash
#

echo "Starting apache httpd"
apachectl -k start -DFOREGROUND
